# crontab
# comment
* * * * * ls -l
first-18 first-last * * mon-thu ls -l
