# crontab
# comment
* * * * * ls -l
first-15 first-last * * thu ls -l
