# crontab
# comment
* * * * * ls -l
first-last first-last * * thu ls -l
